# home/.config/systemd/user/wallthis.timer
[Unit]
Description=Change wallpaper every so often

[Timer]
OnBootSec=1min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=default.target
